# SCRIPTS/safrano9999/image/services/tailscale/tailscale-up.service.d/10-tailscale-ssh.conf
[Service]
ExecStart=
ExecStart=/bin/bash -c 'args=(up --authkey="$$TS_AUTHKEY" --accept-routes --accept-dns --ssh); [ -n "$$TS_HOSTNAME" ] && args+=(--hostname="$$TS_HOSTNAME"); exec /usr/bin/tailscale "$${args[@]}"'
ExecStartPost=-/usr/bin/tailscale set --ssh
