# services/systemd/openclaw.service.d/20-safrano9999.conf
[Unit]
Requires=openclaw-safrano9999.service
After=openclaw-safrano9999.service
