# services/systemd/openclaw-config.service.d/20-codex-models.conf
[Service]
ExecStartPost=/usr/local/bin/openclaw-codex-models
