# services/systemd/openclaw-config.service.d/20-provider-discovery.conf
[Service]
ExecStartPost=/usr/local/bin/openclaw_provider_discovery.py
