# services/openclaw.service
[Unit]
Description=OpenClaw gateway
After=tailscale-up.service openclaw-config.service
Wants=tailscale-up.service

[Service]
Type=simple
PassEnvironment=OPENCLAW_START OPENCLAW_GATEWAY_TOKEN OPENCLAW_GATEWAY_PORT OPENCLAW_TELEGRAM_CHAT_ID OPENCLAW_TELEGRAMTOKEN SAFRANO9999_WEBHOOK_SCRIPT SAFRANO9999_FULLRUN_SCRIPT SAFRANO9999_FULLRUN_DELAY SAFRANO9999_FULLRUN_ON_START OPENCLAW_CRONTABS OPENCLAW_CRONTAB SAFRANO9999_CRONTAB SAFRANO9999_ROUTINES_CRONTAB OPENCLAW_CRON_CHANNEL OPENCLAW_CRON_TO OPENCLAW_CRON_THREAD_ID OPENCLAW_CRON_SETUP_DELAY OPENCLAW_CRON_AGENT OPENCLAW_CRON_SESSION OPENCLAW_CRON_TZ OPENCLAW_CRON_COMMAND
# Fedora OpenAI-v1 provider variables are activated by the image drop-in.
Environment=HOME=/root
ExecCondition=/bin/bash -c '[ "$${OPENCLAW_START:-1}" = "1" ]'
ExecStart=/bin/bash -lc 'args=(gateway run --bind lan --port "$${OPENCLAW_GATEWAY_PORT:-18789}"); if [ -n "$${OPENCLAW_GATEWAY_TOKEN:-}" ]; then args+=(--auth token --token "$$OPENCLAW_GATEWAY_TOKEN"); fi; exec /usr/local/bin/openclaw "$${args[@]}"'
ExecStartPost=-/bin/bash -lc 'if [ -x /usr/local/bin/openclaw-crontabs ]; then (sleep "$${OPENCLAW_CRON_SETUP_DELAY:-20}"; /usr/local/bin/openclaw-crontabs) & fi'
ExecStartPost=-/bin/bash -lc 'if [ "$${SAFRANO9999_FULLRUN_ON_START:-1}" = "1" ] && [ -x "$${SAFRANO9999_FULLRUN_SCRIPT:-/usr/local/bin/safrano9999-fullrun}" ]; then (sleep "$${SAFRANO9999_FULLRUN_DELAY:-100}"; "$${SAFRANO9999_FULLRUN_SCRIPT:-/usr/local/bin/safrano9999-fullrun}") & fi'
Restart=always
TimeoutStartSec=180

[Install]
WantedBy=multi-user.target
