# SCRIPTS/safrano9999/image/services/readme/safrano9999-welcome.service
[Unit]
Description=Build safcontainer welcome documentation
After=fedora44-ai.service
Before=openclaw.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/safrano9999-welcome
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
