# services/systemd/tailscale-up.service.d/10-tailscale-ssh.conf
[Service]
PassEnvironment=TS_EXTRA_ARGS
