# SCRIPTS/safrano9999/image/services/tailscale/tailscaled.service
[Unit]
Description=Tailscale daemon
After=network.target

[Service]
Type=simple
PassEnvironment=TS_AUTHKEY TS_STATE_DIR
Environment=TS_STATE_DIR=/var/lib/tailscale
ExecCondition=/bin/bash -c '[ -n "$$TS_AUTHKEY" ]'
ExecStart=/usr/sbin/tailscaled --state=${TS_STATE_DIR}/tailscaled.state
Restart=on-failure

[Install]
WantedBy=multi-user.target
